# etc/systemd/atlas-bot.service
[Unit]
Description=ATLAS Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Update User and paths below to match your installation
User=atlas
WorkingDirectory=/opt/atlas-bot
ExecStart=/opt/atlas-bot/venv/bin/python bot.py
Restart=always
RestartSec=10
Environment=PATH=/opt/atlas-bot/venv/bin:/usr/local/bin:/usr/bin
EnvironmentFile=/opt/atlas-bot/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas-bot

[Install]
WantedBy=multi-user.target
